# root/etc/services.d/coredns/run
#!/usr/bin/with-contenv bash

cd /home/cuong/config/coredns

if netstat -apn | grep -q ":53 "; then
  echo "Another service is using port 53, disabling CoreDNS"
  sleep infinity
else
  exec \
    /home/cuong/app/coredns -dns.port=53
fi